# apps/reaper/test/integration/supervisord.conf
[supervisord]
nodaemon=true

[program:vault-server]
autorestart=false
command="vault server -dev"
priority=100

[program:vault-setup]
command=/setup-vault.sh
priority=200
